# services/powertop.service
[Unit]
Description=Powertop Auto Tune + Auto USB Mouse Detection
After=multi-user.target systemd-udev-settle.service
Requires=systemd-udev-settle.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/powertop --auto-tune
ExecStartPost=/bin/bash -c 'for d in /sys/bus/usb/devices/*; do [ -e "$d/power/control" ] || continue; prod=$(cat "$d/product" 2>/dev/null | tr "[:upper:]" "[:lower:]"); manu=$(cat "$d/manufacturer" 2>/dev/null | tr "[:upper:]" "[:lower:]"); echo "$prod $manu" | grep -iqE "usb (receiver|pointer|mouse)|mouse|logitech|dell|hp" && echo on > "$d/power/control" || echo auto > "$d/power/control"; done'

[Install]
WantedBy=multi-user.target
